# conf/supervisor.ini
[program:hypertexthero.com]
command=/home/hth/webapps/hypertexthero/hypertexthero-env/bin/python2.7 /home/hth/webapps/hypertexthero/hypertexthero-env/bin/gunicorn_django --daemon -b 127.0.0.1:23035 -w 2 --max-requests 500
directory=/home/hth/webapps/hypertexthero/hypertexthero-env/hth/hth
user=hth
autostart=true
autorestart=true
redirect_stderr=True